# contrib/srvfb.service
[Unit]
Description=Framebuffer Server

[Service]
ExecStart=/usr/bin/srvfb -device /dev/fb0
Restart=always

[Install]
WantedBy=multi-user.target
